# scripts/leaderboard/leaderboard.cron
# /etc/cron.d/leaderboard.cron
* * * * * ec2-user /opt/ioquake3/leaderboard.sh /mnt/server0.log /mnt/server1.log
# You need a newline at the end of your crontab, or cron ignores the file.
