# scripts/stats/stats.cron
# /etc/cron.d/leaderboard.cron
0 */6 * * * ec2-user /opt/ioquake3/stats.sh > /mnt/stats.csv
# You need a newline at the end of your crontab, or cron ignores the file.
